# roles/common/automatic-shutdown/templates/automatic-shutdown.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Automatic shutdown service timer

[Timer]
OnCalendar=Mon..Thu,Sun *-*-* 23:00:00
OnCalendar=Sat,Sun *-*-* 02:00:00

[Install]
WantedBy=timers.target
