# roles/common/media-server/templates/custom-docker-image-updater.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Custom docker image updater timer

[Timer]
OnCalendar=Sat *-*-* 04:00:00

[Install]
WantedBy=timers.target
